# systemd/lora_cam_bridge.service

[Unit]
Description = LoRa-triggered elephant camera
After = multi-user.target

[Service]
ExecStart = /home/a/JumboShoo/venv/bin/python /home/a/JumboShoo/scripts/lora_cam_bridge.py
WorkingDirectory = /home/a/JumboShoo/scripts
User = a
StandardOutput = journal+console
StandardError = inherit
TTYPath = /dev/tty1
Restart = on-failure
RestartSec = 5

[Install]
WantedBy = multi-user.target
