# systemd/reboot.service

[Unit]
Description = Reboot the Pi Regularly

[Service]
Type = oneshot
ExecStart = /sbin/shutdown -r now


